# dfs.service
[Service]
ExecStart=/usr/bin/dfs
Description="This is dezetralized file storage service"
